# connmanupd-network.service
[Unit]
Description=Update connman network information

[Service]
Type=oneshot
ExecStart=/usr/bin/connmanupd -n

[Install]
WantedBy=multi-user.target
